# deploy/neko-sensei.service
[Unit]
Description=Neko Sensei (Next.js)
After=network.target

[Service]
WorkingDirectory=/opt/neko-sensei
EnvironmentFile=-/opt/neko-sensei/.env.local
Environment=NODE_ENV=production
Environment=PORT=3500
Environment=HOSTNAME=127.0.0.1
Environment=NEXT_PUBLIC_BASE_PATH=/neko-sensei
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
User=xno

[Install]
WantedBy=multi-user.target
